# docker/armv7l/stf-armv7l-publish-latest.timer
[Unit]
Description=Run stf-armv7l-publish-latest.service daily

[Timer]
OnCalendar=*-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
